# gratiamonitor/SOURCES/gratiagraph-monitor/gratiagraph-monitor.cron
*/10 * * * * root /usr/sbin/gratiagraph-monitor.sh
