# config/gratia_mon.cron
*/10 * * * * root /sbin/service GratiaWeb restart
