# config/gratia_data.cron
01 0-23/8 * * *   root static_graphs >> /var/log/static_graphs.log 2>>/var/log/static_graphs.err
#01-59/10 * * * *  root gip_record
####01 0-23/8 * * *   root DBPARAM_LOCATION=/etc/DBParam.xml gratia_voms_query
####01 0 * * *        root gip_subcluster_record
####59 0-23/8 * * *   root cms_storage_record
#01 01 1,10,20 * * root wlcg_pledge_email
#01 04 * * *       root ligo_summary
#03 04 * * *       root cms_summary
####20 */2 * * *      root DBPARAM_LOCATION=/etc/DBParam.xml rsv_calc
